# mockingbird.service
[Unit]

After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mockingbird
Restart=always
RestartSec=60
